# profile.d/system-wide-settings.sh
#
# System-wide aliases
#

alias du0='du -hs'
alias du1='du -h --max-depth=1'
alias ll='ls -FGlAh --color=auto --group-directories-first --indicator-style=file-type'
alias ls='ls --color=auto'
alias top='htop'
